# rootfs/etc/services.d/sphere/finish
#!/usr/bin/execlineb -S0

rm /etc/sphere/sphere.pid
s6-svscanctl -t /var/run/s6/services
