# rootfs/etc/services.d/sphere/run
#!/usr/bin/with-contenv bash

cd /etc/sphere && exec ./spheresvr
