# rootfs/etc/services.d/sphere/run
#!/usr/bin/with-contenv bash

cd /etc/sphere
script --return -c "/etc/sphere/spheresvr" /dev/null